# templates/cantaloupe.service.j2
[Unit]
Description=Cantaloupe IIIF Image Server

[Service]
ExecStart={{ cantaloupe_start_command }}
User={{ cantaloupe_user }}

[Install]
WantedBy=multi-user.target
